# DEBIAN/postinst
#!/bin/bash

echo "$0 : $@"

# reloading systemd services
systemctl daemon-reload

systemctl --no-pager enable avahi-backup-client.service
systemctl --no-pager restart avahi-backup-client.service
systemctl --no-pager status avahi-backup-client.service
echo "=========================================================="
echo "* Avahi-Backup Client"
echo "* Please configure /etc/avahi-backup.d for manual sources."
echo "* Please run avahi-backup.sh client-init for client one"
echo "* time setup."
echo "=========================================================="


exit 0